# debian/prerm
#!/bin/sh
# prerm script for emby
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    remove|upgrade|deconfigure)
        . /usr/share/debconf/confmodule
        db_version 2.0
        db_capb backup

        EMBY_PID=`cat /var/run/emby.pid 2> /dev/null || true`
        if [ -n "$EMBY_PID" ]; then
            echo "Stopping Emby Server!"
            if [ -f /etc/init.d/emby ]; then
                    # try graceful termination;
                if which service >/dev/null 2>&1; then
                    service emby stop
                elif which invoke-rc.d >/dev/null 2>&1; then
                    invoke-rc.d emby stop
                else
                    /etc/init.d/emby stop
                fi
            fi
            kill -USR1 $EMBY_PID
            sleep 1
            [ -f /var/run/emby.pid ] && rm /var/run/emby.pid
            if ps -p $EMBY_PID > /dev/null; then
                db_fset emby/old-running seen false || true
                db_input critical emby/old-running || true
                db_go || true
                exit 1
            fi
        fi
        # stop mediabrowser if old binary
        MB_PID=`cat /var/run/mediabrowser.pid 2> /dev/null || true`
        if [ -n "$MB_PID" ]; then
            echo "Stopping MediaBrowser Server!"
            if [ -f /etc/init.d/emby ]; then
                # try graceful termination;
                if which service >/dev/null 2>&1; then
                    service emby stop
                elif which invoke-rc.d >/dev/null 2>&1; then
                    invoke-rc.d emby stop
                else
                    /etc/init.d/emby stop
                fi
            fi
            kill -USR1 $MB_PID
            sleep 1
            [ -f /var/run/emby.pid ] && rm /var/run/emby.pid
            if ps -p $MB_PID > /dev/null; then
                db_fset emby/old-running seen false || true
                db_input critical emby/old-running || true
                db_go || true
                exit 1
            fi
        fi
    ;;

    failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
