# docker-containers/base/overlay-common/etc/services.d/emby-server/run
#!/usr/bin/execlineb -P
# vim:set ft=sh sw=2 sts=2 st=2 et:

with-contenv

exec emby-server start
