# docker-containers/testing/root/etc/services.d/emby-server/run
#!/usr/bin/execlineb -P

exec emby-server start
